# docker/php/crontab
#!/bin/sh

php /var/www/symfony/bin/console swiftmailer:spool:send