# docker/php/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7]
command=php-fpm -F
